# config/table.service
[Unit]
Description=Light-table manager
After=nginx.service

[Service]
ExecStart=/usr/local/bin/table

[Install]
WantedBy=default.target
